# packaging/fossilsafe.service
[Unit]
Description=FossilSafe LTO Tape Backup Appliance
Documentation=https://fossilsafe.io
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
WorkingDirectory={{INSTALL_DIR}}
Environment="PATH={{VENV_DIR}}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH={{INSTALL_DIR}}"
Environment="FOSSILSAFE_CONFIG_PATH={{CONFIG_PATH}}"
Environment="FOSSILSAFE_DATA_DIR={{DATA_DIR}}"
ExecStart={{VENV_DIR}}/bin/gunicorn \
  --config={{INSTALL_DIR}}/gunicorn.conf.py \
  --bind={{BACKEND_BIND}}:{{BACKEND_PORT}} \
  --workers=2 \
  --timeout=120 \
  --access-logfile=- \
  --error-logfile=- \
  backend.lto_backend_main:app

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{DATA_DIR}} {{CONFIG_PATH}}
PrivateTmp=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=no
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=infinity

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fossilsafe

# Device access for tape operations
SupplementaryGroups={{TAPE_GROUP}}

[Install]
WantedBy=multi-user.target