# packaging/fossilsafe.service
[Unit]
Description=FossilSafe LTO Tape Backup Appliance
Documentation=https://fossilsafe.io
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
WorkingDirectory={{INSTALL_DIR}}

# Environment configuration
Environment="PATH={{VENV_DIR}}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=-/etc/default/fossilsafe
EnvironmentFile=-/etc/fossilsafe/config.json

# Core service execution
ExecStart={{VENV_DIR}}/bin/gunicorn \
  --config {{INSTALL_DIR}}/gunicorn.conf.py \
  --bind {{BACKEND_BIND}}:{{BACKEND_PORT}} \
  --workers 4 \
  --worker-class sync \
  --timeout 300 \
  --access-logfile - \
  --error-logfile - \
  backend.lto_backend_main:app

ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s

# Process isolation & security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{DATA_DIR}} {{INSTALL_DIR}}/backend /var/log/fossilsafe /var/www/fossilsafe /etc/fossilsafe

# Resource limits
MemoryLimit=2G
CPUQuota=80%
TasksMax=1024

# Tape device access (udev-managed group)
SupplementaryGroups={{TAPE_GROUP}}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fossilsafe

[Install]
WantedBy=multi-user.target
